# systemd/clash-proxy-refresh.service
[Unit]
Description=Render and reconcile Clash Proxy runtime config
Documentation=file:/root/clash_proxy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/clash_proxy/systemd/clash-proxy-refresh.sh
