# systemd-user/cproxy-refresh.service
[Unit]
Description=Render cproxy runtime config
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cproxy render
