# systemd-user/cproxy-refresh.timer
[Unit]
Description=Periodically render cproxy runtime config

[Timer]
Unit=cproxy-refresh.service
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
